# deploy/supervisord.conf
[supervisord]

[program:omni-comms-worker-1]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan consume:rabbitmq --queue=EmailQueue
stdout_logfile=/var/www/storage/logs/worker.log
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopwaitsecs=3600

[program:omni-comms-worker-2]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=EmailQueue,default
stdout_logfile=/var/www/storage/logs/worker.log
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopwaitsecs=3600

[program:omni-comms-worker-3]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan consume:rabbitmq --queue=SMSQueue
stdout_logfile=/var/www/storage/logs/worker.log
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopwaitsecs=3600

[program:omni-comms-worker-4]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=SMSQueue,default
stdout_logfile=/var/www/storage/logs/worker.log
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopwaitsecs=3600

[program:omni-comms-worker-5]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=default,EmailQueue
stdout_logfile=/var/www/storage/logs/worker.log
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stopwaitsecs=3600